# .config/systemd/user/low-battery.service
[Unit]
Description=Low Battery Warning

[Service]
ExecStart=$HOME/.dotfiles/scripts/low_battery.sh
