# .config/systemd/user/obsidian-daily-commit.service
[Unit]
Description=Daily Obsidian Commit

[Service]
Type=simple
ExecStart=/home/allisnull/.dotfiles/scripts/obsidian_commit.sh
