# .config/systemd/user/obsidian-daily-commit.timer
[Unit]
Description=Hourly Obsidian Vault commit

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
